# deployment/service.service
[Unit]
Description=YoutubeDownloader
After=network.target

[Service]
Type=simple
User=yourusername
WorkingDirectory=/path/to/YoutubeDownloader
ExecStart=/path/to/YoutubeDownloader/.venv/bin/python src/main.py
Restart=always
RestartSec=10

# Environment variables (optional)
# Environment="TASK_RETENTION_MINUTES=60"
# Environment="TASK_CLEANUP_INTERVAL_SECONDS=120"

# Security settings (optional, recommended)
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/path/to/data

[Install]
WantedBy=multi-user.target